# deploy/systemd/anesthesia-exam-telegram-status.service
[Unit]
Description=Anesthesia Exam Telegram Status Report
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=__RUN_USER__
Group=__RUN_GROUP__
WorkingDirectory=__PROJECT_DIR__
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-__PROJECT_DIR__/.env
ExecStart=__PROJECT_DIR__/.venv/bin/python __PROJECT_DIR__/scripts/run_telegram_status_report.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
