# deploy/systemd/anesthesia-exam-openclaw-worker.timer
[Unit]
Description=Run Anesthesia Exam OpenClaw Backlog Worker periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true
Unit=anesthesia-exam-openclaw-worker.service

[Install]
WantedBy=timers.target
