# deploy/systemd/anesthesia-exam-web.service
[Unit]
Description=Anesthesia Exam Streamlit Web
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_GROUP__
WorkingDirectory=__PROJECT_DIR__
Environment=PYTHONUNBUFFERED=1
Environment=EXAM_AGENT_PROVIDER=opencode
Environment=EXAM_OPENCODE_MODEL=gb10/Qwen3.5-122B-A10B-Q5_K_M-00001-of-00003.gguf
ExecStart=/usr/bin/env bash __PROJECT_DIR__/scripts/run_web.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=20
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
