# deploy/systemd/anesthesia-exam-telegram-bot.service
[Unit]
Description=Anesthesia Exam Telegram Admin Bot
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_GROUP__
WorkingDirectory=__PROJECT_DIR__
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-__PROJECT_DIR__/.env
ExecStart=__PROJECT_DIR__/.venv/bin/python __PROJECT_DIR__/scripts/run_telegram_admin_bot.py
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
